# files/systemd/check_kvm.service
[Unit]
Description=Get the latest KVM status information

[Service]
Type=simple
ExecStart=/opt/check_kvm.d/update
