# files/systemd/check_kvm.service
[Unit]
Description=Get the latest KVM status information

[Service]
Type=simple
ExecStart=/usr/lib64/nagios/plugins/check_kvm.d/update
